# iot/sensors/mock/lm35/crontab
*/1 * * * * root /usr/local/bin/node /lm35/lm35.js > /var/log/cron.log 2>&1
